# openstack-ceilometer/openstack-ceilometer-collector.upstart
description "OpenStack Ceilometer collector agent"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/ceilometer-collector --config-file /etc/ceilometer/ceilometer.conf --logfile /var/log/ceilometer/collector.log" ceilometer
