# openstack-ceilometer-slc6/openstack-ceilometer-central.upstart
description "OpenStack Ceilometer central agent"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/ceilometer-agent-central --logfile /var/log/ceilometer/central.log" ceilometer
