# openstack-ceilometer-slc6/openstack-ceilometer-alarm-evaluator.upstart
description "OpenStack Ceilometer alarm evaluation service"

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [S016]

respawn

exec su -s /bin/sh -c "exec /usr/bin/ceilometer-alarm-evaluator --logfile /var/log/ceilometer/alarm-evaluator.log" ceilometer
